# package/base-files/files/etc/uci-defaults/99-wan-dhcp
#!/bin/sh
# CapOS default network: server-style (all NICs = WAN DHCP)

changed=0

# 确保 network 配置存在
[ -f /etc/config/network ] || touch /etc/config/network

# 删除 OpenWrt 路由器默认的 lan / br-lan
uci -q delete network.lan && changed=1
uci -q delete network.globals && changed=1

# 删除 bridge 设备（如果存在）
for dev in $(uci show network | grep "=device" | cut -d. -f2 | cut -d= -f1); do
    name="$(uci -q get network.$dev.name)"
    case "$name" in
        br-*)
            uci delete network.$dev
            changed=1
            ;;
    esac
done

# 为所有物理以太网口创建 DHCP 接口
for iface in $(ip -o link show | awk -F': ' '{print $2}' | grep -E '^e(th|n|m)[0-9]+'); do
    # 跳过已存在的接口
    uci -q get network.$iface >/dev/null && continue

    uci set network.$iface=interface
    uci set network.$iface.device="$iface"
    uci set network.$iface.proto='dhcp'
    changed=1
done

# IPv6：DHCPv6 client（无 PD 假设）
for iface in $(uci show network | grep "=interface" | cut -d. -f2 | cut -d= -f1); do
    proto="$(uci -q get network.$iface.proto)"
    [ "$proto" = "dhcp" ] || continue

    uci set network.${iface}6=interface
    uci set network.${iface}6.device="@${iface}"
    uci set network.${iface}6.proto='dhcpv6'
    changed=1
done

[ "$changed" -eq 1 ] && uci commit network

exit 0
